# templates/00-RequireMountsFor.service.j2
[Unit]
RequireMountsFor={{nginx_RequireMountsFor}}
