# files/letsencrypt-daily.cron
#!/usr/bin/env bash

/usr/bin/certbot renew
